# supervisord_conf.d/appflowy_cloud.conf
[program:appflowy_cloud]
command=appflowy_cloud
autorestart=true
