# supervisord_conf.d/redis.conf
[program:redis]
command=redis-server --bind 0.0.0.0 --protected-mode no
autorestart=true
